# gerrit/files/gerrit.service
[Unit]
Description={{service_name}}

[Service]
ExecStart={{directory}}/bin/gerrit.sh run
User={{user}}